# crontab
TZ=Asia/Tokyo
PATH="/usr/local/bin:$PATH"

0 7 * * * cd /app/ && python paper_summary/main.py > /proc/1/fd/1 2>/proc/1/fd/2
